# install/templates/embedding-model.service
[Unit]
Description=Porrima llama.cpp Embedding Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=${LLAMA_BIN_DIR}
${ENVIRONMENT_LINES}
ExecStart=${LLAMA_SERVER_BIN} \
    -m ${EMBEDDING_MODEL_PATH} \
    --alias ${EMBEDDING_MODEL_ID} \
    --embedding \
    --pooling mean \
    --port 8084 \
    --host 127.0.0.1 \
    --n-gpu-layers 0 \
    --ctx-size 8192 \
    --batch-size 8192 \
    --ubatch-size 8192
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
